# root/etc/systemd/system/archlinux-mirror-rsync.timer
[Unit]
Description=Archlinux mirror server rsync servier timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
