# root/etc/systemd/system/efistub-update.path
[Unit]
Description=Copy EFISTUB Kernel to UEFISYS Partition

[Path]
PathChanged=/boot/

[Install]
WantedBy=multi-user.target
WantedBy=system-update.target
